# server/supervisor.conf
[supervisord] 
logfile = /home/docker/logs/supervisor.log
logfile_maxbytes = 10MB
childlogdir = /home/docker/logs/

[program:nginx]
command = /usr/sbin/nginx